# debian/grase-repo.cron.daily
#!/bin/bash

repoexist2=`grep "packages.grasehotspot.org/ purewhite main" /etc/apt/sources.list.d/* /etc/apt/sources.list`
    if [ "$repoexist2" = "" ]
    then
         touch /etc/apt/sources.list.d/grasehotspot.list
         echo "deb http://packages.grasehotspot.org/ purewhite main" > /etc/apt/sources.list.d/grasehotspot.list
    fi

apt-get update -qq

apt-get install -qq -y `dpkg -l 'grase-*' | grep ii | awk '{ print $2; }'`
exit 0

